# packaging/systemd/cleanup-rbd.service
[Unit]
Description=Cleanup unused RBD mappings

[Service]
Type=oneshot
ExecStart=cleanup_rbd
