# packaging/systemd/cleanup-rbd.service
[Unit]
Description=Cleanup unused RBD mappings

[Service]
Type=oneshot
ExecStart=/usr/bin/cleanup_rbd
